# uol_rpi_tbot/etc/lcas.crontab
*/1 * * * * /home/lcas/workspace/src/network-scripts/iptracker.sh

