# uol_rpi_tbot/etc/supervisor/conf.d/tmule.conf
[program:tmule]
command=/usr/local/bin/tmule server
process_name=%(program_name)s 
directory=/home/lcas
user=lcas
autostart=true
autorestart=true 
redirect_stderr=true 
 
